# make.d/incus/systemd/rke2-network-config.service
[Unit]
Description=RKE2 Network Configuration Service
After=systemd-networkd.service cloud-init.service
Wants=systemd-networkd.service
Before=rke2-install.service

[Service]
Type=oneshot
ExecStart=/.rke2lab/scripts.d/rke2-network-config.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
