# make.d/incus/systemd/rke2-server.service.d/kubeconfig.conf
[Service]
ExecStartPost=/bin/sh -xc '/.rke2lab/scripts.d/rke2-tools-configuration-directories.sh'
ExecStartPost=/bin/sh -xc '/.rke2lab/scripts.d/rke2-vip-kubeconfig.sh'
