# make.d/incus/systemd/rke2-server.service.d/override.conf
[Service]
Restart=always
RestartSec=10s
[Install]
WantedBy=multi-user.target
