# make.d/incus/systemd/rke2-server.service.d/post-start.conf
[Service]
ExecStartPost=/bin/bash -c 'source <(flox activate --dir /var/lib/rancher/rke2) && until kubectl get --raw /readyz &>/dev/null; do echo "Waiting for API server..."; sleep 5; done && echo "API server is ready"'
