# make.d/incus/systemd/rke2-server.service.d/pre-start.conf
[Unit]
Requires=rke2-remount-shared.service rke2-install.service
Wants=rke2-remount-shared.service
After=rke2-remount-shared.service rke2-install.service
[Service]
ExecStartPre=/bin/sh -xc '/.rke2lab/scripts.d/rke2-server-pre-start.sh'
ExecStartPre=/bin/sh -xc '/.rke2lab/scripts.d/rke2-config-install.sh'
ExecStartPost=/bin/sh -xc '/.rke2lab/scripts.d/rke2-server-post-start.sh'
